# raspberry/shutdown/shutdown.service
[Unit]
Description=Shutdown Button Daemon

[Service]
ExecStart =/usr/bin/python3 /opt/shutdown.py

[Install]
WantedBy=multi-user.target
